# deploy/systemd/nova_backtest.service
[Unit]
Description=Trading Nova - backtest
After=network.target

[Service]
ExecStart=/root/tradingbot/bin/python3 -u /root/tradingbot/backtest.py
Restart=always
RestartSec=10
WorkingDirectory=/root/tradingbot

[Install]
WantedBy=multi-user.target
