# deploy/systemd/nova_market_regime.service
[Unit]
Description=Trading Nova - market_regime
After=network.target

[Service]
ExecStart=/root/tradingbot/bin/python3 -u /root/tradingbot/market_regime.py
Restart=always
RestartSec=10
WorkingDirectory=/root/tradingbot

[Install]
WantedBy=multi-user.target
