# deploy/systemd/nova_sentiment_bot.service
[Unit]
Description=Trading Nova - sentiment_bot
After=network.target

[Service]
ExecStart=/root/tradingbot/bin/python3 -u /root/tradingbot/sentiment_bot.py
Restart=always
RestartSec=10
WorkingDirectory=/root/tradingbot

[Install]
WantedBy=multi-user.target
