# deploy/systemd/nova_weekly_digest.service
[Unit]
Description=Trading Nova - weekly_digest
After=network.target

[Service]
ExecStart=/root/tradingbot/bin/python3 -u /root/tradingbot/weekly_digest.py
Restart=always
RestartSec=10
WorkingDirectory=/root/tradingbot

[Install]
WantedBy=multi-user.target
